# services/led-monitor.service
[Unit]
Description=Orange Pi Control Center — LED Status Monitor
After=weather-station.service
Wants=weather-station.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/bash /opt/orangepi/scripts/led_monitor.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orangepi-led

[Install]
WantedBy=multi-user.target
