# services/weather-station.service
[Unit]
Description=Orange Pi Control Center — Flask Applications
Documentation=https://github.com/HicusDanielle/OrangePi
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=kiosk
WorkingDirectory=/opt/orangepi/apps
ExecStart=/usr/bin/bash /opt/orangepi/scripts/start_all.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orangepi-apps

# Resource limits — prevent runaway processes on 1 GB RAM
MemoryMax=512M
TasksMax=64

[Install]
WantedBy=multi-user.target
