# services/x-display.service
[Unit]
Description=Orange Pi Control Center — X Display & Chromium Kiosk
Documentation=https://github.com/HicusDanielle/OrangePi
After=network-online.target weather-station.service
Wants=weather-station.service
BindsTo=weather-station.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=kiosk
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/usr/bin/bash /opt/orangepi/scripts/start_display.sh
Restart=on-failure
RestartSec=8
TimeoutStartSec=90

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orangepi-display

[Install]
WantedBy=multi-user.target
